# deploy/systemd/disaster-alert.service
[Unit]
Description=Disaster Alert Bark subscription service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/disaster-alert
EnvironmentFile=-/opt/disaster-alert/.env
ExecStart=/opt/disaster-alert/disaster-alert-backend
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
